# templates/web-stack.service.j2
[Unit]
Description=Web Stack (Traefik + Nginx) Docker Compose Service
Requires=docker.service
After=docker.service
BindsTo=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ vault_admin_username }}
Group={{ vault_admin_username }}
WorkingDirectory={{ web_stack_dir }}
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=0
Restart=on-failure
RestartSec=5

# Environment variables
Environment=COMPOSE_PROJECT_NAME=webstack
Environment=COMPOSE_FILE=docker-compose.yml

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ web_stack_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=web-stack

[Install]
WantedBy=multi-user.target
